# crontab/restart.sh
#!/bin/bash

# Reiniciar la Raspberry Pi
sudo reboot
